# container-files/etc/supervisor.d/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D
stdout_logfile=/data/logs/sshd.log
stderr_logfile=/data/logs/sshd.log